# package/ubuntu/control/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/mireka" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d mireka stop || true
	else
		/etc/init.d/mireka stop || true
	fi
fi
