# pack/deb/deb/prerm
#!/bin/sh
set -e

systemctl --now disable mireka || true
